# supervisord-mysqld.conf
[program:mysqld]
command=/start.sh
numprocs=1
autostart=true
autorestart=true